# docker/supervisor-conf.d/gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn run:app -b 0.0.0.0:5000 --log-level=debug
directory=/srv/

